# deploy/systemd/user/interest-radar-collect.timer
[Unit]
Description=Run Interest Radar collection daily

[Timer]
OnCalendar=*-*-* 07:40:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
