# deploy/systemd/user/interest-radar-digest.timer
[Unit]
Description=Send Interest Radar digest daily

[Timer]
OnCalendar=*-*-* 08:30:00 Asia/Shanghai
Persistent=true
RandomizedDelaySec=5m

[Install]
WantedBy=timers.target
