# deploy/systemd/user/interest-radar-backup.service
[Unit]
Description=Interest Radar SQLite backup

[Service]
Type=oneshot
WorkingDirectory=%h/meowbot
EnvironmentFile=%h/meowbot/.env
ExecStart=%h/meowbot/bin/radar backup
MemoryMax=96M
CPUQuota=25%
Nice=10
NoNewPrivileges=true

