# deploy/systemd/user/meowbot.service
[Unit]
Description=Meowbot personal secretary
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/meowbot
EnvironmentFile=%h/meowbot/.env
ExecStart=%h/meowbot/bin/meowbot bot
Restart=always
RestartSec=5
MemoryMax=96M
CPUQuota=20%
Nice=10
NoNewPrivileges=true

[Install]
WantedBy=default.target

